# package/etc/init/pf_ring.conf
# pfring - 2003-13 (C) ntop.org

description "PF_RING"
author "Luca Deri <deri@ntop.org> and Alfredo Cardigliano <cardigliano@ntop.org>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Start the process
exec /etc/init.d/pf_ring start

