# package/ubuntu/debian/postinst
#!/bin/sh

/sbin/depmod

echo "Adding the pf_ring startup script"
update-rc.d pf_ring defaults 30 >/dev/null
update-rc.d cluster defaults 40 >/dev/null
/sbin/ldconfig > /dev/null 2>&1

exit 0
