# package/ubuntu/debian/postinst
#!/bin/sh

/sbin/depmod

echo "Adding the pf_ring startup script"
update-rc.d pf_ring defaults 30 >/dev/null
update-rc.d cluster defaults 40 >/dev/null

if [ -f /opt/napatech3/lib/libntapi.so ]; then cp /opt/napatech3/lib/libntapi.so /lib64; fi
if [ -f /opt/napatech3/lib/libntos.so ];  then cp /opt/napatech3/lib/libntos.so /lib64;  fi
if [ -f /opt/snf/lib/libsnf.so ];         then cp /opt/snf/lib/libsnf.so /lib64; ln -s /lib64/libsnf.so /lib64/libsnf.so.0; fi
if [ -f /opt/accolade/lib/libanic.so ];   then cp /opt/accolade/lib/libanic.so /lib64;   fi

/sbin/ldconfig > /dev/null 2>&1

exit 0
